# docker_overlay/overlays/01_core_configuration/overlay/usr/lib/systemd/system/resize_fs.service
[Unit]
Description=Resize Root Partition
ConditionPathExists=/opt/ovos/resize_fs

[Service]
Type=oneshot
ExecStart=/bin/bash /opt/ovos/resize_root_partition.sh

[Install]
WantedBy=multi-user.target