# docker_overlay/overlays/03_sj201/overlay/usr/lib/systemd/system/pulseaudio.service
[Unit]
Description=PulseAudio Sound System
After=syslog.target
After=avahi-daemon.service network.target

[Service]
ExecStart=/usr/bin/pulseaudio --system --disallow-exit --daemonize=no -v
Restart=always

[Install]
WantedBy=multi-user.target

