# docker_overlay/overlays/03_sj201/overlay/usr/lib/systemd/system/sj201-shutdown.service
[Unit]
Description=SJ201 Shutdown Service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStop=/bin/bash /opt/ovos/shutdown_sj201.sh

[Install]
WantedBy=multi-user.target